# supervisor.jenkins-slave.conf
[program:jenkins-slave]
command=/usr/local/bin/jenkins-slave.sh
autorestart=true
user=jenkins
